# deploy/services/datetime.timer
[Unit]
Description = Adjust time

[Timer]
OnBootSec=10
OnUnitActiveSec = 6h
RandomizedDelaySec = 10m

[Install]
WantedBy = timers.target
